# crontab
*/5 * * * * docker exec agorakit.local php artisan schedule:run > /dev/null
*/5 * * * * docker exec db.local /db_maintenance.sh > /dev/null
5 */6 * * * docker exec agorakit.local /backup.sh > /dev/null
