# crontab
*/5 * * * * docker exec agorakit.local php artisan schedule:run
